# crontab
*/3 * * * * /home/lawrence/configRepo/config/scripts/batterynotify

